# templates/systemd/system/clamav-daemon.socket.d/z_custom_socket.conf.j2
# This file is prefixed with `z` because otherwise
# the system-provided `extend.conf` overrides our settings

[Socket]
ListenStream=
ListenStream={{clamd_socket_path}}

